# scripts/gitdash.service
[Unit]
Description=gitdash — local git status dashboard
Documentation=https://github.com/imwebdev/gitdash
After=network.target

# Rate-limit restart loops — these directives belong in [Unit], not [Service]
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
WorkingDirectory={{REPO_DIR}}
ExecStart=/bin/bash {{REPO_DIR}}/bin/gitdash start
EnvironmentFile=%h/.config/gitdash/env

# Restart on crash with a short backoff
Restart=on-failure
RestartSec=5

# Logging via journald — view with: journalctl --user -u gitdash -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
